# systemd/mbsync.timer
[Unit]
Description=Mailbox synchronization timer

[Timer]
OnStartupSec=2m
OnUnitInactiveSec=10m

[Install]
WantedBy=timers.target
